# arms.cron
# /etc/cron.d/arms: rsync abs backups to remote host(s).
#
# From the crontab(5) man page:
#
# field          allowed values
# -----          --------------
# minute         0-59
# hour           0-23
# day of month   1-31
# month          1-12 (or names, see below)
# day of week    0-7 (0 or 7 is Sun, or use names)
#

MAILTO="root"

0  20 *  *  * root /usr/local/sbin/arms mirror_host.mydomain /var/local/backups /var/local/

# /etc/cron.d/arms EOF
